# gb10-gateway/kundali-gb10-gateway.service
[Unit]
Description=Kundali GB10 gateway (thin auth proxy in front of local Ollama)
After=network.target

[Service]
Type=simple
# deploy-local.sh installs the gateway to this runtime dir (decoupled from
# the git checkout so branch switches / worktrees never break the service).
WorkingDirectory=%h/.local/share/kundali-gb10-gateway
ExecStart=%h/.local/share/kundali-gb10-gateway/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8100
Restart=on-failure
RestartSec=5
# Secrets + Supabase coordinates live in an untracked env file created by
# deploy-local.sh (GB10_INTERNAL_SECRET, GB10_PAID_SECRET, SUPABASE_URL,
# SUPABASE_ANON_KEY, GB10_ALLOWED_MODELS, OLLAMA_URL).
EnvironmentFile=%h/.local/share/kundali-gb10-gateway/gateway.env

[Install]
WantedBy=default.target
